# oled_Nanopi_Neo2/oled.service
[Unit]
Description=MoOde Audio OLED Display
After=multi-user.target
 
[Service]
Type=simple
ExecStart=/usr/bin/python /root/oled_Nanopi_Neo2/oled.py
Restart=on-failure
 
[Install]
WantedBy=multi-user.target
